# keychain/zsh/profile.d/keychain.sh
# Sets up keychain agents on login.

$HOME/util/keychain/ssh-agent.sh
source $HOME/.ssh/agent_info > /dev/null

# TODO: set up gpg-agent?
